# installer/transact-qda.service
[Unit]
Description=transact-qda
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=/usr/bin/env node __INSTALL_DIR__/build/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
